# raspberry/etc/cron.d/icpa
01 0     * * *     root   /icpa/bin/save_data

